# infra/systemd/wx-ybhz-cavp.service
[Unit]
Description=宜泊慧智 C-AVP 后端 (wx_ybhzCavpSystem, HTTP 12380)
After=network.target
Wants=network-online.target

[Service]
Type=simple
# 部署前须存在：
#   useradd --system --home /opt/cavp --shell /usr/sbin/nologin cavp
#   mkdir -p /opt/cavp/wx_ybhzCavpSystem /etc/cavp
# 下列路径按现网实际调整（本机常用 /home/wx_ybhzCavpSystem）
User=root
WorkingDirectory=/home/wx_ybhzCavpSystem
EnvironmentFile=-/home/wx_ybhzCavpSystem/cavp.env
# 服务器上 git pull 后执行 mvn -DskipTests package，产物在 target/
# JAVA_HOME / java 路径按机器调整（which java）
# 本机 HTTP :12380；对外 HTTPS 由 Nginx 等反代
ExecStart=/root/miniconda3/bin/java -jar /home/wx_ybhzCavpSystem/target/wx-ybhz-cavp-system-1.0.0.jar --spring.profiles.active=prod
Restart=on-failure
RestartSec=5

# 勿默认开 PrivateTmp：部分云主机/内核会报 226/NAMESPACE

[Install]
WantedBy=multi-user.target
