# infra/systemd/wx-ybhz-cavp.service
[Unit]
Description=宜泊慧智 C-AVP 后端 (wx_ybhzCavpSystem, HTTPS 9065)
After=network.target
Wants=network-online.target

[Service]
Type=simple
# 部署前须存在：
#   useradd --system --home /opt/cavp --shell /usr/sbin/nologin cavp
#   mkdir -p /opt/cavp/wx_ybhzCavpSystem /opt/cavp/cert /etc/cavp
User=cavp
Group=cavp
WorkingDirectory=/opt/cavp/wx_ybhzCavpSystem
EnvironmentFile=-/etc/cavp/cavp.env
# 服务器上 git pull 后执行 mvn -DskipTests package，产物在 target/
# JAVA_HOME / java 路径按机器调整（which java）
ExecStart=/usr/bin/java -jar /opt/cavp/wx_ybhzCavpSystem/target/wx-ybhz-cavp-system-1.0.0.jar --spring.profiles.active=prod
Restart=on-failure
RestartSec=5

# 勿默认开 PrivateTmp：部分云主机/内核会报 226/NAMESPACE

[Install]
WantedBy=multi-user.target
